# templates/kibana.legacy.service.lsbheader.erb
#!/bin/sh
#
# kibana    Web interface to Logstash.
#
# chkconfig:   345 99 01
# description: Kibana is a web interface to Logstash.

### BEGIN INIT INFO
# Provides: kibana
# Required-Start: $all
# Required-Stop: $all
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Kibana
# Description: Kibana is a web interface to Logstash.
### END INIT INFO

DESC="Kibana 4"
NAME=kibana

